# docker/crontab
# crontab 文件
0 9 * * * cd /app && python main.py >> /var/log/cron.log 2>&1
# 空行是必须的
